# .config/autostart/emacs-daemon.desktop
[Desktop Entry]
Type=Application
Name=Emacs Daemon
Exec=emacs --daemon
Comment=Start Emacs as a daemon
